# packages/fedora/src/main/scripts/fedora
#!/bin/bash
#
# fedora        
#
# chkconfig: 345 63 37
# description:  Start up the fedora app.

# Source function library.
. /etc/init.d/functions

# Start multicast for ecache
/usr/local/topaz/bin/multicast

FEDORA_HOME="/usr/local/topaz/fedora-0.5-SNAPSHOT"
export FEDORA_HOME

FEDORA_SH="$FEDORA_HOME/server/bin/fedora.sh"
if [ ! -f "$FEDORA_SH" ] ; then
    echo "Fedora script not found: $FEDORA_SH"
    exit 1
fi

RETVAL=0

case "$1" in
    start)
        echo "Starting Fedora"
        /bin/su topaz -s /bin/sh -c "\"$FEDORA_SH\" start mysql \"$@\""
        RETVAL=$?
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/fedora
        ;;

    stop)
        echo "Stopping Fedora"
        /bin/su topaz -s /bin/sh -c "\"$FEDORA_SH\" stop mysql \"$@\""
        RETVAL=$?
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/fedora
        ;;

    restart)
        /bin/su topaz -s /bin/sh -c "\"$FEDORA_SH\" restart \"$@\""
        ;;

    status)
        /bin/su topaz -s /bin/sh -c "\"$FEDORA_SH\" status \"$@\""
        RETVAL=$?
        ;;

    *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

exit $RETVAL
