# packages/topaz/src/main/scripts/postinstall
# Register the topaz service
/sbin/chkconfig --add topaz
