# packages/topaz/src/main/scripts/postinstall
# Register the topaz service
/bin/cp /etc/topaz/rc.d/topaz /etc/init.d/topaz
/sbin/chkconfig --add topaz

