# deploy/etc/init.d/multicast
#!/usr/bin/env bash
# $HeadURL:: http://gandalf/svn/head/deploy/etc/init.d/mulgara                                    $
# $Id: mulgara 1343 2006-11-26 06:49:27Z ebrown $
#
# Tomcat init script
# Copyright (c) 2006 by Topaz, Inc.
# http://topazproject.org
#
# Licensed under the Educational Community License version 1.0
# http://opensource.org/licenses/ecl1.php

# chkconfig: 345 92 14
# description: ecache multicast configuration

ifconfig eth0 multicast
ifconfig lo multicast
route add -net 224.0.0.0 netmask 240.0.0.0 dev eth0
route add -net 224.0.0.0 netmask 240.0.0.0 dev lo

