# packages/topaz-fedora/src/main/scripts/preremove
if [ $1 = 0 ]; then
    service topaz_fedora stop > /dev/null 2>&1
    /sbin/chkconfig --del topaz_fedora || :
fi
