# packages/cas/src/main/scripts/postinstall
# Register the cas service
/bin/cp /etc/topaz/rc.d/cas /etc/init.d/cas
/sbin/chkconfig --add cas

