# packages/cas/src/main/scripts/postinstall
# Register the cas service
/sbin/chkconfig --add cas
