# packages/plosone/src/main/scripts/postinstall
# Register the PLoS ONE application
/sbin/chkconfig --add plosone
