# packages/plosone/src/main/scripts/postinstall
# Register the PLoS ONE application
/bin/cp /etc/topaz/rc.d/plosone /etc/init.d/plosone
/sbin/chkconfig --add plosone

