# packages/mulgara/src/main/scripts/postinstall
# Register the mulgara service
/sbin/chkconfig --add mulgara
