# packages/mulgara/src/main/scripts/postinstall
# $HeadURL::                                                                   $
# $Id$

# Register the mulgara service
/sbin/chkconfig --add mulgara

# Create symbolic link in /usr/local/topaz
PKG=mulgara
LINK=/usr/local/topaz/$PKG
rm -f $LINK
DIR=`ls -td /usr/local/topaz/*"$PKG"-* 2>/dev/null | head -1`
[ -n "$DIR" ] && ln -s $DIR $LINK
exit 0
