# packages/mulgara/src/main/scripts/postinstall
# Register the mulgara service
/bin/cp /etc/topaz/rc.d/mulgara /etc/init.d/mulgara
/sbin/chkconfig --add mulgara

